# linux/supervisord.conf
; RepoFabric supervisord configuration.
; Three programs: node admin (browser-facing admin UI on :8086 plus the
; Express installer server on :8091), pwsh bridge (compute + worker pool),
; and cron (catalog walker + scheduled sync enqueue).
;
; Setup mode: entrypoint.sh writes /var/lib/repofabric/setup-mode if the config
; file is missing. While that flag exists, cron stays stopped and only the
; /setup/* routes in node admin are active. After the wizard completes,
; the node admin's setup handler removes the flag and supervisorctl is
; signalled to start cron.

[supervisord]
nodaemon=true
user=root
logfile=/var/lib/repofabric/logs/supervisord.log
pidfile=/var/lib/repofabric/supervisord.pid
childlogdir=/var/lib/repofabric/logs

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ---- node admin ------------------------------------------------------------
; Always running. While setup-mode flag is present, server.js mounts the
; /setup/* router and refuses every other route.
[program:node-admin]
command=/usr/bin/node /opt/repofabric-admin/src/server.js
user=repofabric
environment=NODE_ENV="production",PORT="8086"
autostart=true
autorestart=true
startsecs=3
startretries=10
stdout_logfile=/var/lib/repofabric/logs/node-admin.out.log
stderr_logfile=/var/lib/repofabric/logs/node-admin.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; ---- pwsh bridge -----------------------------------------------------------
; Long-running pwsh process. Owns the SQLite handle, the worker pool, and
; the loopback HTTP listener on 127.0.0.1:8085.
[program:pwsh-bridge]
command=/usr/bin/pwsh -NoLogo -NoProfile -Command "Import-Module /opt/repofabric/src/RepoFabric.psd1; Start-RfWebUI -ListenPrefix http://127.0.0.1:8085/"
user=repofabric
autostart=true
autorestart=true
startsecs=5
startretries=10
stdout_logfile=/var/lib/repofabric/logs/pwsh-bridge.out.log
stderr_logfile=/var/lib/repofabric/logs/pwsh-bridge.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; ---- cron ------------------------------------------------------------------
; Walker (every 5 min) and scheduled-sync enqueue (every 6 hours by default,
; overridable via service.yaml). The cron daemon itself auto-starts; the
; individual crontab entries gate themselves on setup.complete so they
; noop until the wizard finishes. A previous gating approach kept the cron
; PROGRAM down until setup.js sent a 'start repofabric:cron' supervisorctl
; signal, but that left cron STOPPED across every container recreate
; because setup.complete already existed and setup.js never ran again.
[program:cron]
command=/usr/sbin/cron -f
user=root
autostart=true
autorestart=true
startsecs=2
startretries=5
stdout_logfile=/var/lib/repofabric/logs/cron.out.log
stderr_logfile=/var/lib/repofabric/logs/cron.err.log

; ---- ConfigFabric pwsh bridge (absorbed sidecar) ---------------------------
; The co-hosted ConfigFabric PowerShell module + loopback listener on
; 127.0.0.1:8089 (distinct from RepoFabric's 8085). autostart=false: the
; entrypoint runs `supervisorctl start repofabric:cf-pwsh-bridge` ONLY when
; CONFIGFABRIC_ENABLED=true, so a standalone RepoFabric never starts it. Owns
; ConfigFabric's own SQLite at /var/lib/configfabric/state.sqlite.
[program:cf-pwsh-bridge]
command=/usr/bin/pwsh -NoLogo -NoProfile -Command "Import-Module /opt/configfabric/src/ConfigFabric.psd1; Start-CfWebUI -ListenPrefix http://127.0.0.1:8089/"
user=repofabric
environment=CONFIGFABRIC_STATE_DIR="/var/lib/configfabric",CONFIGFABRIC_PUBLISHER_URL="http://127.0.0.1:8089"
autostart=false
autorestart=true
startsecs=5
startretries=10
stdout_logfile=/var/lib/repofabric/logs/cf-pwsh-bridge.out.log
stderr_logfile=/var/lib/repofabric/logs/cf-pwsh-bridge.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; ---- initial upstream index seed (one-shot) --------------------------------
; Seeds the upstream winget index ONCE, right after first-run setup completes,
; so a freshly stood-up instance has a searchable package index immediately
; rather than waiting up to 6h for the first scheduled sync. The script waits
; for setup.complete, runs Update-RfUpstreamIndex once, and writes a marker so
; it never repeats. autorestart=false: it exits after seeding; if it errors, the
; next scheduled sync refreshes the index. Progress shows in the admin UI via
; the same index-refresh status the manual "Refresh index" button uses.
[program:initial-index]
command=/usr/local/bin/initial-index.sh
user=repofabric
autostart=true
autorestart=false
startsecs=0
startretries=0
stdout_logfile=/var/lib/repofabric/logs/initial-index.out.log
stderr_logfile=/var/lib/repofabric/logs/initial-index.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

[group:repofabric]
programs=node-admin,pwsh-bridge,cron,cf-pwsh-bridge,initial-index
